# deploy/systemd/aether-daemon.service
# Hardened systemd unit for running the Aether daemon directly on the host
# (an alternative to the Docker Compose 'daemon' service — running on the host
# avoids giving a container the Docker socket).
#
# Install:
#   sudo cp deploy/systemd/aether-daemon.service /etc/systemd/system/
#   sudo useradd --system --no-create-home --shell /usr/sbin/nologin aether || true
#   sudo usermod -aG docker aether
#   sudo mkdir -p /var/lib/aether/volumes /var/lib/aether/backups
#   sudo chown -R aether:aether /var/lib/aether /opt/aether
#   sudo systemctl daemon-reload && sudo systemctl enable --now aether-daemon

[Unit]
Description=Aether game-server daemon
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=aether
Group=aether
WorkingDirectory=/opt/aether/apps/daemon
ExecStart=/usr/bin/node /opt/aether/apps/daemon/dist/index.js
Restart=always
RestartSec=3

# Environment (or use EnvironmentFile=/opt/aether/.env)
Environment=NODE_ENV=production
Environment=DAEMON_PORT=8080
Environment=DAEMON_DATA_DIR=/var/lib/aether/volumes
Environment=DAEMON_BACKUP_DIR=/var/lib/aether/backups
Environment=DOCKER_SOCKET=/var/run/docker.sock
# EnvironmentFile=/opt/aether/.env

# ── Hardening ──────────────────────────────────────────────
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/aether
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
SystemCallArchitectures=native
CapabilityBoundingSet=

[Install]
WantedBy=multi-user.target
